# docker/crontab
* * * * * cd /var/app && flask training_session:run > /dev/stdout
